# zmoviemagic.service
[Unit]
Description=zmoviemagic Fire TV Cinema web app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=afradkin
WorkingDirectory=/home/afradkin/zmoviemagic
ExecStart=/home/afradkin/zmoviemagic/venv/bin/python app.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
