# zmoviemagic.service
[Unit]
Description=ZMovieMagic - Disney+ Fire TV Controller
After=network.target

[Service]
# User running the application (change 'ubuntu' to your actual username if different)
User=ubuntu
Group=ubuntu

# Path to the application directory
WorkingDirectory=/home/ubuntu/zmoviemagic

# Environment variables
# Make sure .env exists in the WorkingDirectory
Environment="PATH=/home/ubuntu/zmoviemagic/venv/bin"

# Command to start the application using gunicorn
# Ensure 'gunicorn' is installed in your venv
ExecStart=/home/ubuntu/zmoviemagic/venv/bin/gunicorn --workers 3 --bind 0.0.0.0:8000 app:app

# Restart automatically if it crashes
Restart=always

[Install]
WantedBy=multi-user.target
