# zmoviemagic.service
[Unit]
Description=ZMovieMagic - Disney+ Fire TV Controller
After=network.target

[Service]
# User running the application (change 'ubuntu' to your actual username if different)
User=afradkin
Group=afradkin

# Path to the application directory
WorkingDirectory=/home/afradkin/zmoviemagic

# Environment variables
# Make sure .env exists in the WorkingDirectory
Environment="PATH=/home/afradkin/zmoviemagic/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Command to start the application using gunicorn
# Ensure 'gunicorn' is installed in your venv
ExecStart=/home/afradkin/zmoviemagic/venv/bin/gunicorn --workers 3 --bind 0.0.0.0:8000 app:app

# Restart automatically if it crashes
Restart=always

[Install]
WantedBy=multi-user.target
